# bin/local.tifmac.kasa-backup.sample.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>local.tifmac.kasa-backup.svn</string>
	<key>Disabled</key>
	<false/>
	<key>LowPriorityIO</key>
	<true/>
	<key>Program</key>
	<string>/Users/tufanoruk/bin/kasa-backup.sh</string>
	<key>ProgramArguments</key>
	<array>
		<string>kasa-backup.sh</string>
		<string>/usr/local/svn</string>
		<string>tifmac.local.tr-svn</string>
	</array>
	<key>WatchPaths</key>
    <array>
        <string>/usr/local/svn</string>
    </array>
    <key>ThrottleInterval</key>
    <integer>3600</integer>
</dict>
</plist>
